# supervisord-v2.conf
[supervisord]
nodaemon=true
logfile=/var/log/supervisord.log
pidfile=/var/run/supervisord.pid

[program:affinity]
command=/app/start-mcp.sh affinity 5001
environment=PORT=5001
directory=/app/mcp_servers/affinity
autostart=true
autorestart=true
stderr_logfile=/var/log/affinity.err.log
stdout_logfile=/var/log/affinity.out.log
stopasgroup=true
killasgroup=true

[program:attio]
command=/app/start-mcp.sh attio 5002
environment=PORT=5002
directory=/app/mcp_servers/attio
autostart=true
autorestart=true
stderr_logfile=/var/log/attio.err.log
stdout_logfile=/var/log/attio.out.log
stopasgroup=true
killasgroup=true

[program:close]
command=/app/start-mcp.sh close 5003
environment=PORT=5003
directory=/app/mcp_servers/close
autostart=true
autorestart=true
stderr_logfile=/var/log/close.err.log
stdout_logfile=/var/log/close.out.log
stopasgroup=true
killasgroup=true

[program:hubspot]
command=/app/start-mcp.sh hubspot 5004
environment=PORT=5004
directory=/app/mcp_servers/hubspot
autostart=true
autorestart=true
stderr_logfile=/var/log/hubspot.err.log
stdout_logfile=/var/log/hubspot.out.log
stopasgroup=true
killasgroup=true

[program:salesforce]
command=/app/start-mcp.sh salesforce 5005
environment=PORT=5005
directory=/app/mcp_servers/salesforce
autostart=true
autorestart=true
stderr_logfile=/var/log/salesforce.err.log
stdout_logfile=/var/log/salesforce.out.log
stopasgroup=true
killasgroup=true

[program:gong]
command=/app/start-mcp.sh gong 5006
environment=PORT=5006
directory=/app/mcp_servers/gong
autostart=true
autorestart=true
stderr_logfile=/var/log/gong.err.log
stdout_logfile=/var/log/gong.out.log
stopasgroup=true
killasgroup=true

[program:asana]
command=/app/start-mcp.sh asana 5007
environment=PORT=5007
directory=/app/mcp_servers/asana
autostart=true
autorestart=true
stderr_logfile=/var/log/asana.err.log
stdout_logfile=/var/log/asana.out.log
stopasgroup=true
killasgroup=true

[program:clickup]
command=/app/start-mcp.sh clickup 5008
environment=PORT=5008
directory=/app/mcp_servers/clickup
autostart=true
autorestart=true
stderr_logfile=/var/log/clickup.err.log
stdout_logfile=/var/log/clickup.out.log
stopasgroup=true
killasgroup=true

[program:jira]
command=/app/start-mcp.sh jira 5009
environment=PORT=5009
directory=/app/mcp_servers/jira
autostart=true
autorestart=true
stderr_logfile=/var/log/jira.err.log
stdout_logfile=/var/log/jira.out.log
stopasgroup=true
killasgroup=true

[program:linear]
command=/app/start-mcp.sh linear 5010
environment=PORT=5010
directory=/app/mcp_servers/linear
autostart=true
autorestart=true
stderr_logfile=/var/log/linear.err.log
stdout_logfile=/var/log/linear.out.log
stopasgroup=true
killasgroup=true

[program:monday]
command=/app/start-mcp.sh monday 5011
environment=PORT=5011
directory=/app/mcp_servers/monday
autostart=true
autorestart=true
stderr_logfile=/var/log/monday.err.log
stdout_logfile=/var/log/monday.out.log
stopasgroup=true
killasgroup=true

[program:motion]
command=/app/start-mcp.sh motion 5012
environment=PORT=5012
directory=/app/mcp_servers/motion
autostart=true
autorestart=true
stderr_logfile=/var/log/motion.err.log
stdout_logfile=/var/log/motion.out.log
stopasgroup=true
killasgroup=true

[program:freshdesk]
command=/app/start-mcp.sh freshdesk 5013
environment=PORT=5013
directory=/app/mcp_servers/freshdesk
autostart=true
autorestart=true
stderr_logfile=/var/log/freshdesk.err.log
stdout_logfile=/var/log/freshdesk.out.log
stopasgroup=true
killasgroup=true

[program:mailchimp]
command=/app/start-mcp.sh mailchimp 5014
environment=PORT=5014
directory=/app/mcp_servers/mailchimp
autostart=true
autorestart=true
stderr_logfile=/var/log/mailchimp.err.log
stdout_logfile=/var/log/mailchimp.out.log
stopasgroup=true
killasgroup=true

[program:resend]
command=/app/start-mcp.sh resend 5015
environment=PORT=5015
directory=/app/mcp_servers/resend
autostart=true
autorestart=true
stderr_logfile=/var/log/resend.err.log
stdout_logfile=/var/log/resend.out.log
stopasgroup=true
killasgroup=true

[program:slack]
command=/app/start-mcp.sh slack 5016
environment=PORT=5016
directory=/app/mcp_servers/slack
autostart=true
autorestart=true
stderr_logfile=/var/log/slack.err.log
stdout_logfile=/var/log/slack.out.log
stopasgroup=true
killasgroup=true

[program:discord]
command=/app/start-mcp.sh discord 5017
environment=PORT=5017
directory=/app/mcp_servers/discord
autostart=true
autorestart=true
stderr_logfile=/var/log/discord.err.log
stdout_logfile=/var/log/discord.out.log
stopasgroup=true
killasgroup=true

[program:gmail]
command=/app/start-mcp.sh gmail 5018
environment=PORT=5018
directory=/app/mcp_servers/gmail
autostart=true
autorestart=true
stderr_logfile=/var/log/gmail.err.log
stdout_logfile=/var/log/gmail.out.log
stopasgroup=true
killasgroup=true

[program:google_calendar]
command=/app/start-mcp.sh google_calendar 5019
environment=PORT=5019
directory=/app/mcp_servers/google_calendar
autostart=true
autorestart=true
stderr_logfile=/var/log/google_calendar.err.log
stdout_logfile=/var/log/google_calendar.out.log
stopasgroup=true
killasgroup=true

[program:google_meet]
command=/app/start-mcp.sh google_meet 5020
environment=PORT=5020
directory=/app/mcp_servers/google_meet
autostart=true
autorestart=true
stderr_logfile=/var/log/google_meet.err.log
stdout_logfile=/var/log/google_meet.out.log
stopasgroup=true
killasgroup=true

[program:twilio]
command=/app/start-mcp.sh twilio 5021
environment=PORT=5021
directory=/app/mcp_servers/twilio
autostart=true
autorestart=true
stderr_logfile=/var/log/twilio.err.log
stdout_logfile=/var/log/twilio.out.log
stopasgroup=true
killasgroup=true

[program:whatsapp]
command=/app/start-mcp.sh whatsapp 5022
environment=PORT=5022
directory=/app/mcp_servers/whatsapp
autostart=true
autorestart=true
stderr_logfile=/var/log/whatsapp.err.log
stdout_logfile=/var/log/whatsapp.out.log
stopasgroup=true
killasgroup=true

[program:cal_com]
command=/app/start-mcp.sh cal_com 5023
environment=PORT=5023
directory=/app/mcp_servers/cal_com
autostart=true
autorestart=true
stderr_logfile=/var/log/cal_com.err.log
stdout_logfile=/var/log/cal_com.out.log
stopasgroup=true
killasgroup=true

[program:confluence]
command=/app/start-mcp.sh confluence 5024
environment=PORT=5024
directory=/app/mcp_servers/confluence
autostart=true
autorestart=true
stderr_logfile=/var/log/confluence.err.log
stdout_logfile=/var/log/confluence.out.log
stopasgroup=true
killasgroup=true

[program:notion]
command=/app/start-mcp.sh notion 5025
environment=PORT=5025
directory=/app/mcp_servers/notion
autostart=true
autorestart=true
stderr_logfile=/var/log/notion.err.log
stdout_logfile=/var/log/notion.out.log
stopasgroup=true
killasgroup=true

[program:google_docs]
command=/app/start-mcp.sh google_docs 5026
environment=PORT=5026
directory=/app/mcp_servers/google_docs
autostart=true
autorestart=true
stderr_logfile=/var/log/google_docs.err.log
stdout_logfile=/var/log/google_docs.out.log
stopasgroup=true
killasgroup=true

[program:google_sheets]
command=/app/start-mcp.sh google_sheets 5027
environment=PORT=5027
directory=/app/mcp_servers/google_sheets
autostart=true
autorestart=true
stderr_logfile=/var/log/google_sheets.err.log
stdout_logfile=/var/log/google_sheets.out.log
stopasgroup=true
killasgroup=true

[program:google_slides]
command=/app/start-mcp.sh google_slides 5028
environment=PORT=5028
directory=/app/mcp_servers/google_slides
autostart=true
autorestart=true
stderr_logfile=/var/log/google_slides.err.log
stdout_logfile=/var/log/google_slides.out.log
stopasgroup=true
killasgroup=true

[program:google_drive]
command=/app/start-mcp.sh google_drive 5029
environment=PORT=5029
directory=/app/mcp_servers/google_drive
autostart=true
autorestart=true
stderr_logfile=/var/log/google_drive.err.log
stdout_logfile=/var/log/google_drive.out.log
stopasgroup=true
killasgroup=true

[program:dropbox]
command=/app/start-mcp.sh dropbox 5030
environment=PORT=5030
directory=/app/mcp_servers/dropbox
autostart=true
autorestart=true
stderr_logfile=/var/log/dropbox.err.log
stdout_logfile=/var/log/dropbox.out.log
stopasgroup=true
killasgroup=true

[program:onedrive]
command=/app/start-mcp.sh onedrive 5031
environment=PORT=5031
directory=/app/mcp_servers/onedrive
autostart=true
autorestart=true
stderr_logfile=/var/log/onedrive.err.log
stdout_logfile=/var/log/onedrive.out.log
stopasgroup=true
killasgroup=true

[program:wordpress]
command=/app/start-mcp.sh wordpress 5032
environment=PORT=5032
directory=/app/mcp_servers/wordpress
autostart=true
autorestart=true
stderr_logfile=/var/log/wordpress.err.log
stdout_logfile=/var/log/wordpress.out.log
stopasgroup=true
killasgroup=true

[program:postgres]
command=/app/start-mcp.sh postgres 5033
environment=PORT=5033
directory=/app/mcp_servers/postgres
autostart=true
autorestart=true
stderr_logfile=/var/log/postgres.err.log
stdout_logfile=/var/log/postgres.out.log
stopasgroup=true
killasgroup=true

[program:supabase]
command=/app/start-mcp.sh supabase 5034
environment=PORT=5034
directory=/app/mcp_servers/supabase
autostart=true
autorestart=true
stderr_logfile=/var/log/supabase.err.log
stdout_logfile=/var/log/supabase.out.log
stopasgroup=true
killasgroup=true

[program:airtable]
command=/app/start-mcp.sh airtable 5035
environment=PORT=5035
directory=/app/mcp_servers/airtable
autostart=true
autorestart=true
stderr_logfile=/var/log/airtable.err.log
stdout_logfile=/var/log/airtable.out.log
stopasgroup=true
killasgroup=true

[program:shopify]
command=/app/start-mcp.sh shopify 5036
environment=PORT=5036
directory=/app/mcp_servers/shopify
autostart=true
autorestart=true
stderr_logfile=/var/log/shopify.err.log
stdout_logfile=/var/log/shopify.out.log
stopasgroup=true
killasgroup=true

[program:stripe]
command=/app/start-mcp.sh stripe 5037
environment=PORT=5037
directory=/app/mcp_servers/stripe
autostart=true
autorestart=true
stderr_logfile=/var/log/stripe.err.log
stdout_logfile=/var/log/stripe.out.log
stopasgroup=true
killasgroup=true

[program:quickbooks]
command=/app/start-mcp.sh quickbooks 5038
environment=PORT=5038
directory=/app/mcp_servers/quickbooks
autostart=true
autorestart=true
stderr_logfile=/var/log/quickbooks.err.log
stdout_logfile=/var/log/quickbooks.out.log
stopasgroup=true
killasgroup=true

[program:moneybird]
command=/app/start-mcp.sh moneybird 5039
environment=PORT=5039
directory=/app/mcp_servers/moneybird
autostart=true
autorestart=true
stderr_logfile=/var/log/moneybird.err.log
stdout_logfile=/var/log/moneybird.out.log
stopasgroup=true
killasgroup=true

[program:coinbase]
command=/app/start-mcp.sh coinbase 5040
environment=PORT=5040
directory=/app/mcp_servers/coinbase
autostart=true
autorestart=true
stderr_logfile=/var/log/coinbase.err.log
stdout_logfile=/var/log/coinbase.out.log
stopasgroup=true
killasgroup=true

[program:github]
command=/app/start-mcp.sh github 5041
environment=PORT=5041
directory=/app/mcp_servers/github
autostart=true
autorestart=true
stderr_logfile=/var/log/github.err.log
stdout_logfile=/var/log/github.out.log
stopasgroup=true
killasgroup=true

[program:brave_search]
command=/app/start-mcp.sh brave_search 5042
environment=PORT=5042
directory=/app/mcp_servers/brave_search
autostart=true
autorestart=true
stderr_logfile=/var/log/brave_search.err.log
stdout_logfile=/var/log/brave_search.out.log
stopasgroup=true
killasgroup=true

[program:exa]
command=/app/start-mcp.sh exa 5043
environment=PORT=5043
directory=/app/mcp_servers/exa
autostart=true
autorestart=true
stderr_logfile=/var/log/exa.err.log
stdout_logfile=/var/log/exa.out.log
stopasgroup=true
killasgroup=true

[program:tavily]
command=/app/start-mcp.sh tavily 5044
environment=PORT=5044
directory=/app/mcp_servers/tavily
autostart=true
autorestart=true
stderr_logfile=/var/log/tavily.err.log
stdout_logfile=/var/log/tavily.out.log
stopasgroup=true
killasgroup=true

[program:firecrawl]
command=/app/start-mcp.sh firecrawl 5045
environment=PORT=5045
directory=/app/mcp_servers/firecrawl
autostart=true
autorestart=true
stderr_logfile=/var/log/firecrawl.err.log
stdout_logfile=/var/log/firecrawl.out.log
stopasgroup=true
killasgroup=true

[program:firecrawl_deep_research]
command=/app/start-mcp.sh firecrawl_deep_research 5046
environment=PORT=5046
directory=/app/mcp_servers/firecrawl_deep_research
autostart=true
autorestart=true
stderr_logfile=/var/log/firecrawl_deep_research.err.log
stdout_logfile=/var/log/firecrawl_deep_research.out.log
stopasgroup=true
killasgroup=true

[program:hacker_news]
command=/app/start-mcp.sh hacker_news 5047
environment=PORT=5047
directory=/app/mcp_servers/hacker_news
autostart=true
autorestart=true
stderr_logfile=/var/log/hacker_news.err.log
stdout_logfile=/var/log/hacker_news.out.log
stopasgroup=true
killasgroup=true

[program:linkedin]
command=/app/start-mcp.sh linkedin 5048
environment=PORT=5048
directory=/app/mcp_servers/linkedin
autostart=true
autorestart=true
stderr_logfile=/var/log/linkedin.err.log
stdout_logfile=/var/log/linkedin.out.log
stopasgroup=true
killasgroup=true

[program:youtube]
command=/app/start-mcp.sh youtube 5049
environment=PORT=5049
directory=/app/mcp_servers/youtube
autostart=true
autorestart=true
stderr_logfile=/var/log/youtube.err.log
stdout_logfile=/var/log/youtube.out.log
stopasgroup=true
killasgroup=true

[program:spotify]
command=/app/start-mcp.sh spotify 5050
environment=PORT=5050
directory=/app/mcp_servers/spotify
autostart=true
autorestart=true
stderr_logfile=/var/log/spotify.err.log
stdout_logfile=/var/log/spotify.out.log
stopasgroup=true
killasgroup=true

[program:openrouter]
command=/app/start-mcp.sh openrouter 5051
environment=PORT=5051
directory=/app/mcp_servers/openrouter
autostart=true
autorestart=true
stderr_logfile=/var/log/openrouter.err.log
stdout_logfile=/var/log/openrouter.out.log
stopasgroup=true
killasgroup=true

[program:heygen]
command=/app/start-mcp.sh heygen 5052
environment=PORT=5052
directory=/app/mcp_servers/heygen
autostart=true
autorestart=true
stderr_logfile=/var/log/heygen.err.log
stdout_logfile=/var/log/heygen.out.log
stopasgroup=true
killasgroup=true

[program:markitdown]
command=/app/start-mcp.sh markitdown 5053
environment=PORT=5053
directory=/app/mcp_servers/markitdown
autostart=true
autorestart=true
stderr_logfile=/var/log/markitdown.err.log
stdout_logfile=/var/log/markitdown.out.log
stopasgroup=true
killasgroup=true

[program:pandoc]
command=/app/start-mcp.sh pandoc 5054
environment=PORT=5054
directory=/app/mcp_servers/pandoc
autostart=true
autorestart=true
stderr_logfile=/var/log/pandoc.err.log
stdout_logfile=/var/log/pandoc.out.log
stopasgroup=true
killasgroup=true

[program:report_generation]
command=/app/start-mcp.sh report_generation 5055
environment=PORT=5055
directory=/app/mcp_servers/report_generation
autostart=true
autorestart=true
stderr_logfile=/var/log/report_generation.err.log
stdout_logfile=/var/log/report_generation.out.log
stopasgroup=true
killasgroup=true

[program:mem0]
command=/app/start-mcp.sh mem0 5056
environment=PORT=5056
directory=/app/mcp_servers/mem0
autostart=true
autorestart=true
stderr_logfile=/var/log/mem0.err.log
stdout_logfile=/var/log/mem0.out.log
stopasgroup=true
killasgroup=true

[program:google_jobs]
command=/app/start-mcp.sh google_jobs 5057
environment=PORT=5057
directory=/app/mcp_servers/google_jobs
autostart=true
autorestart=true
stderr_logfile=/var/log/google_jobs.err.log
stdout_logfile=/var/log/google_jobs.out.log
stopasgroup=true
killasgroup=true

[program:mixpanel]
command=/app/start-mcp.sh mixpanel 5058
environment=PORT=5058
directory=/app/mcp_servers/mixpanel
autostart=true
autorestart=true
stderr_logfile=/var/log/mixpanel.err.log
stdout_logfile=/var/log/mixpanel.out.log
stopasgroup=true
killasgroup=true

[program:nginx]
command=nginx -g 'daemon off;'
autostart=true
autorestart=true
stderr_logfile=/var/log/nginx.err.log
stdout_logfile=/var/log/nginx.out.log

