# supervisord.conf
[supervisord]
nodaemon=true
logfile=/var/log/supervisord.log
pidfile=/var/run/supervisord.pid

[program:affinity]
command=node /app/mcp_servers/affinity/index.js
environment=PORT=5001
autostart=true
autorestart=true
stderr_logfile=/var/log/affinity.err.log
stdout_logfile=/var/log/affinity.out.log

[program:attio]
command=node /app/mcp_servers/attio/index.js
environment=PORT=5002
autostart=true
autorestart=true
stderr_logfile=/var/log/attio.err.log
stdout_logfile=/var/log/attio.out.log

[program:close]
command=node /app/mcp_servers/close/index.js
environment=PORT=5003
autostart=true
autorestart=true
stderr_logfile=/var/log/close.err.log
stdout_logfile=/var/log/close.out.log

[program:hubspot]
command=node /app/mcp_servers/hubspot/index.js
environment=PORT=5004
autostart=true
autorestart=true
stderr_logfile=/var/log/hubspot.err.log
stdout_logfile=/var/log/hubspot.out.log

[program:salesforce]
command=node /app/mcp_servers/salesforce/index.js
environment=PORT=5005
autostart=true
autorestart=true
stderr_logfile=/var/log/salesforce.err.log
stdout_logfile=/var/log/salesforce.out.log

[program:gong]
command=node /app/mcp_servers/gong/index.js
environment=PORT=5006
autostart=true
autorestart=true
stderr_logfile=/var/log/gong.err.log
stdout_logfile=/var/log/gong.out.log

[program:asana]
command=node /app/mcp_servers/asana/index.js
environment=PORT=5007
autostart=true
autorestart=true
stderr_logfile=/var/log/asana.err.log
stdout_logfile=/var/log/asana.out.log

[program:clickup]
command=node /app/mcp_servers/clickup/index.js
environment=PORT=5008
autostart=true
autorestart=true
stderr_logfile=/var/log/clickup.err.log
stdout_logfile=/var/log/clickup.out.log

[program:jira]
command=node /app/mcp_servers/jira/index.js
environment=PORT=5009
autostart=true
autorestart=true
stderr_logfile=/var/log/jira.err.log
stdout_logfile=/var/log/jira.out.log

[program:linear]
command=node /app/mcp_servers/linear/index.js
environment=PORT=5010
autostart=true
autorestart=true
stderr_logfile=/var/log/linear.err.log
stdout_logfile=/var/log/linear.out.log

[program:monday]
command=node /app/mcp_servers/monday/index.js
environment=PORT=5011
autostart=true
autorestart=true
stderr_logfile=/var/log/monday.err.log
stdout_logfile=/var/log/monday.out.log

[program:motion]
command=node /app/mcp_servers/motion/index.js
environment=PORT=5012
autostart=true
autorestart=true
stderr_logfile=/var/log/motion.err.log
stdout_logfile=/var/log/motion.out.log

[program:freshdesk]
command=node /app/mcp_servers/freshdesk/index.js
environment=PORT=5013
autostart=true
autorestart=true
stderr_logfile=/var/log/freshdesk.err.log
stdout_logfile=/var/log/freshdesk.out.log

[program:mailchimp]
command=node /app/mcp_servers/mailchimp/index.js
environment=PORT=5014
autostart=true
autorestart=true
stderr_logfile=/var/log/mailchimp.err.log
stdout_logfile=/var/log/mailchimp.out.log

[program:resend]
command=node /app/mcp_servers/resend/index.js
environment=PORT=5015
autostart=true
autorestart=true
stderr_logfile=/var/log/resend.err.log
stdout_logfile=/var/log/resend.out.log

[program:slack]
command=node /app/mcp_servers/slack/index.js
environment=PORT=5016
autostart=true
autorestart=true
stderr_logfile=/var/log/slack.err.log
stdout_logfile=/var/log/slack.out.log

[program:discord]
command=node /app/mcp_servers/discord/index.js
environment=PORT=5017
autostart=true
autorestart=true
stderr_logfile=/var/log/discord.err.log
stdout_logfile=/var/log/discord.out.log

[program:gmail]
command=node /app/mcp_servers/gmail/index.js
environment=PORT=5018
autostart=true
autorestart=true
stderr_logfile=/var/log/gmail.err.log
stdout_logfile=/var/log/gmail.out.log

[program:google_calendar]
command=node /app/mcp_servers/google_calendar/index.js
environment=PORT=5019
autostart=true
autorestart=true
stderr_logfile=/var/log/google_calendar.err.log
stdout_logfile=/var/log/google_calendar.out.log

[program:google_meet]
command=node /app/mcp_servers/google_meet/index.js
environment=PORT=5020
autostart=true
autorestart=true
stderr_logfile=/var/log/google_meet.err.log
stdout_logfile=/var/log/google_meet.out.log

[program:twilio]
command=node /app/mcp_servers/twilio/index.js
environment=PORT=5021
autostart=true
autorestart=true
stderr_logfile=/var/log/twilio.err.log
stdout_logfile=/var/log/twilio.out.log

[program:whatsapp]
command=node /app/mcp_servers/whatsapp/index.js
environment=PORT=5022
autostart=true
autorestart=true
stderr_logfile=/var/log/whatsapp.err.log
stdout_logfile=/var/log/whatsapp.out.log

[program:cal_com]
command=node /app/mcp_servers/cal_com/index.js
environment=PORT=5023
autostart=true
autorestart=true
stderr_logfile=/var/log/cal_com.err.log
stdout_logfile=/var/log/cal_com.out.log

[program:confluence]
command=node /app/mcp_servers/confluence/index.js
environment=PORT=5024
autostart=true
autorestart=true
stderr_logfile=/var/log/confluence.err.log
stdout_logfile=/var/log/confluence.out.log

[program:notion]
command=node /app/mcp_servers/notion/index.js
environment=PORT=5025
autostart=true
autorestart=true
stderr_logfile=/var/log/notion.err.log
stdout_logfile=/var/log/notion.out.log

[program:google_docs]
command=node /app/mcp_servers/google_docs/index.js
environment=PORT=5026
autostart=true
autorestart=true
stderr_logfile=/var/log/google_docs.err.log
stdout_logfile=/var/log/google_docs.out.log

[program:google_sheets]
command=node /app/mcp_servers/google_sheets/index.js
environment=PORT=5027
autostart=true
autorestart=true
stderr_logfile=/var/log/google_sheets.err.log
stdout_logfile=/var/log/google_sheets.out.log

[program:google_slides]
command=node /app/mcp_servers/google_slides/index.js
environment=PORT=5028
autostart=true
autorestart=true
stderr_logfile=/var/log/google_slides.err.log
stdout_logfile=/var/log/google_slides.out.log

[program:google_drive]
command=node /app/mcp_servers/google_drive/index.js
environment=PORT=5029
autostart=true
autorestart=true
stderr_logfile=/var/log/google_drive.err.log
stdout_logfile=/var/log/google_drive.out.log

[program:dropbox]
command=node /app/mcp_servers/dropbox/index.js
environment=PORT=5030
autostart=true
autorestart=true
stderr_logfile=/var/log/dropbox.err.log
stdout_logfile=/var/log/dropbox.out.log

[program:onedrive]
command=node /app/mcp_servers/onedrive/index.js
environment=PORT=5031
autostart=true
autorestart=true
stderr_logfile=/var/log/onedrive.err.log
stdout_logfile=/var/log/onedrive.out.log

[program:wordpress]
command=node /app/mcp_servers/wordpress/index.js
environment=PORT=5032
autostart=true
autorestart=true
stderr_logfile=/var/log/wordpress.err.log
stdout_logfile=/var/log/wordpress.out.log

[program:postgres]
command=node /app/mcp_servers/postgres/index.js
environment=PORT=5033
autostart=true
autorestart=true
stderr_logfile=/var/log/postgres.err.log
stdout_logfile=/var/log/postgres.out.log

[program:supabase]
command=node /app/mcp_servers/supabase/index.js
environment=PORT=5034
autostart=true
autorestart=true
stderr_logfile=/var/log/supabase.err.log
stdout_logfile=/var/log/supabase.out.log

[program:airtable]
command=node /app/mcp_servers/airtable/index.js
environment=PORT=5035
autostart=true
autorestart=true
stderr_logfile=/var/log/airtable.err.log
stdout_logfile=/var/log/airtable.out.log

[program:shopify]
command=node /app/mcp_servers/shopify/index.js
environment=PORT=5036
autostart=true
autorestart=true
stderr_logfile=/var/log/shopify.err.log
stdout_logfile=/var/log/shopify.out.log

[program:stripe]
command=node /app/mcp_servers/stripe/index.js
environment=PORT=5037
autostart=true
autorestart=true
stderr_logfile=/var/log/stripe.err.log
stdout_logfile=/var/log/stripe.out.log

[program:quickbooks]
command=node /app/mcp_servers/quickbooks/index.js
environment=PORT=5038
autostart=true
autorestart=true
stderr_logfile=/var/log/quickbooks.err.log
stdout_logfile=/var/log/quickbooks.out.log

[program:moneybird]
command=node /app/mcp_servers/moneybird/index.js
environment=PORT=5039
autostart=true
autorestart=true
stderr_logfile=/var/log/moneybird.err.log
stdout_logfile=/var/log/moneybird.out.log

[program:coinbase]
command=node /app/mcp_servers/coinbase/index.js
environment=PORT=5040
autostart=true
autorestart=true
stderr_logfile=/var/log/coinbase.err.log
stdout_logfile=/var/log/coinbase.out.log

[program:github]
command=node /app/mcp_servers/github/index.js
environment=PORT=5041
autostart=true
autorestart=true
stderr_logfile=/var/log/github.err.log
stdout_logfile=/var/log/github.out.log

[program:brave_search]
command=node /app/mcp_servers/brave_search/index.js
environment=PORT=5042
autostart=true
autorestart=true
stderr_logfile=/var/log/brave_search.err.log
stdout_logfile=/var/log/brave_search.out.log

[program:exa]
command=node /app/mcp_servers/exa/index.js
environment=PORT=5043
autostart=true
autorestart=true
stderr_logfile=/var/log/exa.err.log
stdout_logfile=/var/log/exa.out.log

[program:tavily]
command=node /app/mcp_servers/tavily/index.js
environment=PORT=5044
autostart=true
autorestart=true
stderr_logfile=/var/log/tavily.err.log
stdout_logfile=/var/log/tavily.out.log

[program:firecrawl]
command=node /app/mcp_servers/firecrawl/index.js
environment=PORT=5045
autostart=true
autorestart=true
stderr_logfile=/var/log/firecrawl.err.log
stdout_logfile=/var/log/firecrawl.out.log

[program:firecrawl_deep_research]
command=node /app/mcp_servers/firecrawl_deep_research/index.js
environment=PORT=5046
autostart=true
autorestart=true
stderr_logfile=/var/log/firecrawl_deep_research.err.log
stdout_logfile=/var/log/firecrawl_deep_research.out.log

[program:hacker_news]
command=node /app/mcp_servers/hacker_news/index.js
environment=PORT=5047
autostart=true
autorestart=true
stderr_logfile=/var/log/hacker_news.err.log
stdout_logfile=/var/log/hacker_news.out.log

[program:linkedin]
command=node /app/mcp_servers/linkedin/index.js
environment=PORT=5048
autostart=true
autorestart=true
stderr_logfile=/var/log/linkedin.err.log
stdout_logfile=/var/log/linkedin.out.log

[program:youtube]
command=node /app/mcp_servers/youtube/index.js
environment=PORT=5049
autostart=true
autorestart=true
stderr_logfile=/var/log/youtube.err.log
stdout_logfile=/var/log/youtube.out.log

[program:spotify]
command=node /app/mcp_servers/spotify/index.js
environment=PORT=5050
autostart=true
autorestart=true
stderr_logfile=/var/log/spotify.err.log
stdout_logfile=/var/log/spotify.out.log

[program:openrouter]
command=node /app/mcp_servers/openrouter/index.js
environment=PORT=5051
autostart=true
autorestart=true
stderr_logfile=/var/log/openrouter.err.log
stdout_logfile=/var/log/openrouter.out.log

[program:heygen]
command=node /app/mcp_servers/heygen/index.js
environment=PORT=5052
autostart=true
autorestart=true
stderr_logfile=/var/log/heygen.err.log
stdout_logfile=/var/log/heygen.out.log

[program:markitdown]
command=node /app/mcp_servers/markitdown/index.js
environment=PORT=5053
autostart=true
autorestart=true
stderr_logfile=/var/log/markitdown.err.log
stdout_logfile=/var/log/markitdown.out.log

[program:pandoc]
command=node /app/mcp_servers/pandoc/index.js
environment=PORT=5054
autostart=true
autorestart=true
stderr_logfile=/var/log/pandoc.err.log
stdout_logfile=/var/log/pandoc.out.log

[program:report_generation]
command=node /app/mcp_servers/report_generation/index.js
environment=PORT=5055
autostart=true
autorestart=true
stderr_logfile=/var/log/report_generation.err.log
stdout_logfile=/var/log/report_generation.out.log

[program:mem0]
command=node /app/mcp_servers/mem0/index.js
environment=PORT=5056
autostart=true
autorestart=true
stderr_logfile=/var/log/mem0.err.log
stdout_logfile=/var/log/mem0.out.log

[program:google_jobs]
command=node /app/mcp_servers/google_jobs/index.js
environment=PORT=5057
autostart=true
autorestart=true
stderr_logfile=/var/log/google_jobs.err.log
stdout_logfile=/var/log/google_jobs.out.log

[program:mixpanel]
command=node /app/mcp_servers/mixpanel/index.js
environment=PORT=5058
autostart=true
autorestart=true
stderr_logfile=/var/log/mixpanel.err.log
stdout_logfile=/var/log/mixpanel.out.log

[program:nginx]
command=nginx -g 'daemon off;'
autostart=true
autorestart=true
stderr_logfile=/var/log/nginx.err.log
stdout_logfile=/var/log/nginx.out.log
